# ansible/roles/motd/files/update-motd.d/92-recent-bans
#!/bin/sh

out=/var/run/recent_bans

if [ -f "$out" ]; then
    echo
    cat "$out"
fi
