# ansible/roles/motd/files/update-motd.d/99-tail
#!/bin/sh

if [ -f "/etc/motd.tail" ]; then
    echo
    cat "/etc/motd.tail"
fi
